# Agent/sffd-sync/crontab
*/5 * * * * root python /app/sffd_sync.py >> /proc/1/fd/1 2>&1
